# FlightTrackerPreview.service
# Live web preview of the LED panel via RGBMatrixEmulator (preview.py).
# Install:
#   sudo cp FlightTrackerPreview.service /etc/systemd/system/
#   sudo systemctl daemon-reload && sudo systemctl enable --now FlightTrackerPreview.service
# Serves the emulator browser on the port in emulator_config.json (default 8888).
[Unit]
Description=FlightTracker Live Preview (RGBMatrixEmulator browser mirror of the LED panel)
After=network.target FlightTracker.service

[Service]
Type=simple
ExecStart=/home/pi/FlightTracker/env/bin/python /home/pi/FlightTracker/preview.py
WorkingDirectory=/home/pi/FlightTracker
StandardOutput=append:/home/pi/preview.log
StandardError=append:/home/pi/preview.log
Restart=on-failure
RestartSec=5
User=pi
Environment=PYTHONUNBUFFERED=1
Environment=NO_COLOR=1

[Install]
WantedBy=multi-user.target
